# cluster/nopipe-cluster.service
[Unit]
Description=Nopipe Cluster RPC Server
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/opt/nopipe
EnvironmentFile=/opt/nopipe/.env
ExecStart=/opt/nopipe/bin/nopipe-cluster
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nopipe-cluster

[Install]
WantedBy=multi-user.target
