# cluster/polyclaw-cluster.service
[Unit]
Description=Polyclaw Cluster RPC Server
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/polyclaw
EnvironmentFile=/home/ubuntu/polyclaw/cluster/.env
ExecStart=/home/ubuntu/polyclaw/target/release/polyclaw-cluster
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
